# storm-supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:supervisor]
command=storm supervisor

[program:logviewer]
command=storm logviewer